# native/linux/clear-vis.service
[Unit]
Description=Clear Spotify Visualizer Audio Bridge
Documentation=https://github.com/wktkow/clear-spotify-client
After=pulseaudio.service pipewire-pulse.service

[Service]
Type=simple
ExecStart=%h/.local/bin/vis-capture
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
